# configs/automatic-backup/gitlab-deployment-backup.service
[Unit]
Description=GitLab Deployment Backup

[Service]
ExecStart=/usr/bin/gitlab-deployment-backup.sh
